# docker/scheduler/crontab
# GDPR Scheduler - Runs daily at 2:00 AM
0 2 * * * node ace gdpr:scheduler

# Subscription Renewal - Runs daily at 0:05 AM
5 0 * * * node ace subscription:renew

# Credit Requests Cleanup - Runs weekly on Sunday at 3:00 AM
# Deletes processed credit requests older than 90 days (pending requests are preserved)
0 3 * * 0 node ace cleanup:credit-requests

# Notifications Cleanup - Runs daily at 3:00 AM
# Deletes read notifications older than 30 days
0 3 * * * node ace cleanup:notifications

# Auto-Refill Check - Runs daily at 18:00 (6 PM)
# Warns owners 24h before if there are insufficient credits for scheduled auto-refills
0 18 * * * node ace check:auto-refill
